# ansible/roles/cyhy_dashboard/templates/ncats-webd.service
[Unit]
Description=ncats-webd service

[Service]
Type=simple
Restart=always
RestartSec=1
User=cyhy
ExecStart=/usr/local/bin/ncats-webd -ls production
